# crontab
# m h  dom mon dow   command
26 1 * * * /usr/local/bin/certbot.ash
# An empty line is required at the end of this file for a valid cron file.
